# HelperService/LaunchDaemons/io.nekohasekai.sfavt.helper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>io.nekohasekai.sfavt.helper</string>
	<key>BundleProgram</key>
	<string>Contents/Helpers/RootHelper</string>
	<key>MachServices</key>
	<dict>
		<key>287TTNZF8L.io.nekohasekai.sfavt.helper</key>
		<true/>
	</dict>
	<key>AssociatedBundleIdentifiers</key>
	<array>
		<string>io.nekohasekai.sfavt.standalone</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
</dict>
</plist>
